# installer/templates/multi-portal.service
[Unit]
Description=Multi-Portal (Provedor — Portal Admin + Site)
After=network.target postgresql.service

[Service]
Type=simple
User=root
WorkingDirectory=__INSTALL_DIR__
EnvironmentFile=__INSTALL_DIR__/.env
ExecStart=/usr/bin/node __INSTALL_DIR__/dist/server.js
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
